# storm-deb-packaging/debian/apache-storm.storm-nimbus.upstart
description "storm nimbus service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid storm
setgid storm
chdir /usr/lib/storm

script
    [ -r /etc/default/storm ] && . /etc/default/storm
    [ -r /etc/default/storm-nimbus ] && . /etc/default/storm-nimbus

    exec /usr/bin/storm nimbus
end script
