# storm-deb-packaging/debian/apache-storm.storm-logviewer.upstart
description "storm logviewer service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid storm
setgid storm
chdir /usr/lib/storm

script
    [ -r /etc/default/storm ] && . /etc/default/storm
    [ -r /etc/default/storm-logviewer ] && . /etc/default/storm-logviewer

    exec /usr/bin/storm logviewer
end script
