# storm-deb-packaging/debian/apache-storm.storm-supervisor.upstart
description "storm supervisor service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid storm
setgid storm
chdir /usr/lib/storm

script
    [ -r /etc/default/storm ] && . /etc/default/storm
    [ -r /etc/default/storm-supervisor ] && . /etc/default/storm-supervisor

    exec /usr/bin/storm supervisor
end script
