# storm-deb-packaging/debian/apache-storm.storm-drpc.upstart
description "storm drpc service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid storm
setgid storm
chdir /usr/lib/storm

script
    [ -r /etc/default/storm ] && . /etc/default/storm
    [ -r /etc/default/storm-drpc ] && . /etc/default/storm-drpc

    exec /usr/bin/storm drpc
end script
