# storm-deb-packaging/debian/storm-ui.upstart
description "storm ui service"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid storm
setgid storm
chdir /usr/lib/storm

script
    [ -r /etc/default/storm ] && . /etc/default/storm
    [ -r /etc/default/storm-ui ] && . /etc/default/storm-ui

    exec /usr/bin/storm ui
end script
