# apache-storm/debian/storm-drpc.service
# /etc/systemd/system/storm-drpc.service

# TODO: handle ulimit (in /etc/security/limits.conf?)
# increase the amount of open files, because storm can open lots of those.
# and default value was too small.
# ulimit -n 15000
#
# TODO: check if JVM limits could be set in a new configuration file (storm_env.ini)
# STORM_NIMBUS_JVM_MEMORY="512m"
# STORM_SUPERVISOR_JVM_MEMORY="1024m"
# STORM_UI_JVM_MEMORY="256m"

[Unit]
Description="Apache Storm DRPC Service"

[Service]
ExecStart=/usr/bin/storm drpc
User=storm
Restart=always
PIDFile=/var/run/storm-drpc.pid
